# data/postsrsd.service.in
# PostSRSd systemd service file
# Copyright 2022 Timo Röhling
# SPDX-License-Identifier: FSFUL
#
# The copyright holder gives unlimited permission to copy, distribute and modify
# this file.
[Unit]
Description=Sender Rewriting Scheme daemon for Postfix
Before=postfix.service
After=network.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_SBINDIR@/postsrsd -c @CMAKE_INSTALL_FULL_SYSCONFDIR@/@PROJECT_NAME@.conf

[Install]
WantedBy=multi-user.target
