# data/postsrsd.service.in
# PostSRSd systemd service file
# Copyright 2022-2026 Timo Röhling
# SPDX-License-Identifier: FSFUL
#
# The copyright holder gives unlimited permission to copy, distribute and modify
# this file.

[Unit]
Description=Sender Rewriting Scheme daemon for Postfix
Before=postfix.service
After=network.target

[Service]
Type=notify-reload
ExecStart=@CMAKE_INSTALL_FULL_SBINDIR@/postsrsd -C @POSTSRSD_CONFIGDIR@/@PROJECT_NAME@.conf
NotifyAccess=main
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/spool/postfix @POSTSRSD_CHROOTDIR@
PrivateTmp=yes
PrivateDevices=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictSUIDSGID=yes
SystemCallArchitectures=native
LockPersonality=yes
ProtectKernelTunables=yes
RestrictRealtime=yes
ProtectClock=yes
RestrictNamespaces=yes
CapabilityBoundingSet=~CAP_SYS_PTRACE
ProtectHome=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes

[Install]
WantedBy=multi-user.target
